# deployment/supervisord.conf
[supervisord]

[program:process1]
command=/root/.local/share/apibara/bin/apibara run --allow-env=.env src/strkfarm/harvests.ts --sink-id=hemant_120 --status-server-address=0.0.0.0:4140
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/harvests.log
stderr_logfile=/var/log/supervisor/harvests_err.log

[program:process2]
command=/root/.local/share/apibara/bin/apibara run --allow-env=.env src/strkfarm/deposits-withdraws.ts --sink-id=heamnt_150 --status-server-address=0.0.0.0:4130
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/dep-withdraw.log
stderr_logfile=/var/log/supervisor/dep-withdraw_err.log

[program:process3]
command=/root/.local/share/apibara/bin/apibara run --allow-env=.env src/ekubo/position-fees-collected.ts --sink-id=hemant_170 --status-server-address=0.0.0.0:4120
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/position-fees-collected.log
stderr_logfile=/var/log/supervisor/position-fees-collected_err.log

[program:process4]
command=/root/.local/share/apibara/bin/apibara run --allow-env=.env src/ekubo/position-updated.ts --sink-id=heamnt_190 --status-server-address=0.0.0.0:4110
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/position-updated.log
stderr_logfile=/var/log/supervisor/position-updated_err.log